# contrib/systemd/osmo-bsc-nat.service
[Unit]
Description=Osmocom BSC NAT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/osmo-bsc-nat -c /etc/osmocom/osmo-bsc-nat.cfg
RestartSec=2

[Install]
WantedBy=multi-user.target
